# infra/ansible/roles/aio-cooler/files/hosts/dragon/dynamic_led.service
[Unit]
Description=Dynamic LED based on liquid temperature
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/env python3 /usr/local/bin/dynamic_led.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
