# packages/hardware/.config/systemd/user/battery-monitor.service
[Unit]
Description=Dragon Battery Monitor Check
After=graphical-session.target

[Service]
Type=oneshot
Environment=PATH=%h/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/bin
Environment=XDG_RUNTIME_DIR=%t
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=%t/bus
ExecStart=%h/dotfiles/scripts/hardware/battery-monitor