# packages/hardware/.config/systemd/user/battery-monitor.service
[Unit]
Description=Periodically check battery level

[Service]
Type=oneshot
ExecStart=%h/.local/bin/battery-monitor
