# hosts/goldendragon/dotfiles/.config/systemd/user/fprintd-watchdog.service
[Unit]
Description=Fprintd Watchdog - Restart if device claims are stuck

[Service]
Type=oneshot
ExecStart=%h/.local/bin/fprintd-watchdog

[Install]
WantedBy=default.target
