# hosts/goldendragon/etc/systemd/user/fprintd-watchdog.service
[Unit]
Description=Fprintd Watchdog - Restart if device claims are stuck
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=/home/daofficialwizard/.local/bin/fprintd-watchdog

[Install]
WantedBy=default.target
