# hosts/dragon/liquidctl-dragon.service
[Unit]
Description=Dragon Custom Liquidctl AIO Cooler Configuration
Documentation=https://github.com/liquidctl/liquidctl
After=multi-user.target
Wants=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
Group=root

# Initialize and configure Corsair H100i Platinum SE for Dragon
ExecStart=/bin/bash -c '\
    /usr/bin/liquidctl initialize --match "h100i" && \
    /usr/bin/liquidctl set fan speed --match "h100i" \
        20 20 \
        30 30 \
        40 40 \
        50 60 \
        60 80 \
        70 90 \
        80 100 && \
    /usr/bin/liquidctl set led color fixed 000000 --match "h100i" && \
    /usr/bin/liquidctl set pump speed 70 --match "h100i"'

# Restart the service if it fails
Restart=on-failure
RestartSec=30

# Timeout settings
TimeoutStartSec=60
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/dev
DevicePolicy=closed
DeviceAllow=/dev/hidraw* rw
DeviceAllow=/dev/usb* rw

[Install]
WantedBy=multi-user.target 