# hosts/dragon/liquidctl-dragon.service
[Unit]
Description=Dragon Custom Liquidctl AIO Cooler Configuration
Documentation=https://github.com/liquidctl/liquidctl
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes

# Initialize and configure Corsair H100i Platinum SE for Dragon
ExecStart=/bin/bash -c '\
    /usr/bin/liquidctl initialize --match "h100i" && \
    /usr/bin/liquidctl set fan speed --match "h100i" \
        30 20 \
        40 30 \
        50 50 \
        60 70 \
        70 100 && \
    /usr/bin/liquidctl set led color fixed 0080ff --match "h100i"'

[Install]
WantedBy=multi-user.target
