# containers/openfortivpn/podman/openfortivpn.container
[Container]
Image=ghcr.io/daofficialwizard/openfortivpn:latest
Network=host
Privileged=true
PID=host
AddDevice=/dev/ppp
Volume=./config/openfortivpn.conf:/etc/openfortivpn/config:ro
Volume=/etc/resolv.conf:/host/etc/resolv.conf
Volume=openfortivpn-run:/run/openfortivpn
Environment=VPN_HOST=vpn.avular.com
Environment=VPN_PORT=443
Environment=VPN_DOMAIN=avular.dev
Environment=VPN_DNS_SERVERS=10.10.100.50,10.10.100.11
Environment=SAML_PORT=8020
Environment=DNS_METHOD=auto
Environment=AUTO_CONNECT=false
Environment=CONFIG_PATH=/etc/openfortivpn/config
Environment=API_TCP_PORT=8080
Environment=API_UNIX_SOCKET=/run/openfortivpn/api.sock
Environment=LOG_LEVEL=info
HealthCmd=curl -fsS http://127.0.0.1:8080/health
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=10s

[Service]
Restart=unless-stopped

[Install]
WantedBy=default.target
